# client/lwnsim.run/build/lwnsim.service
# Contents of /lib/systemd/system/rftasmd.service
[Unit]
Description=RF Tasmota to FHEM
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/opt/rftasmd/mqtt-fhem
WorkingDirectory=/opt/rftasmd/

[Install]
WantedBy=multi-user.target
