# client/lwnsim.run/build/lwnsim.service
# Contents of /lib/systemd/system/lwnsim.service
[Unit]
Description=LoraWan Simulator Device Configurator
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=10
ExecStart=/opt/lwnsim/lwnsim
WorkingDirectory=/opt/lwnsim/

[Install]
WantedBy=multi-user.target
